# teamspeak.conf
[program:teamspeak]
user = nobody
command = /home/nobody/teamspeak.sh