# teamspeak.conf
[program:teamspeak]
user = nobody
command = /usr/bin/teamspeak3-server default_voice_port=9987 filetransfer_port=30033 query_port=10011 logpath=/config/ dbsqlpath=/usr/share/teamspeak3-server/sql/