# teamspeak.conf
[program:teamspeak]
user = nobody
command = /usr/bin/teamspeak.sh