# setup/teamspeak.conf
[program:teamspeak]
user = nobody
command = /home/nobody/setup.sh